# deploy/launchd/com.openclaw.tokenranger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.tokenranger</string>
    <key>ProgramArguments</key>
    <array>
        <string>~/.openclaw/services/tokenranger/venv/bin/uvicorn</string>
        <string>main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8100</string>
        <string>--workers</string>
        <string>1</string>
    </array>
    <key>WorkingDirectory</key>
    <string>~/.openclaw/services/tokenranger</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>~/.openclaw/logs/tokenranger.log</string>
    <key>StandardErrorPath</key>
    <string>~/.openclaw/logs/tokenranger.log</string>
</dict>
</plist>
